# debian/vyatta-pre-boot.target
[Unit]
Description=Pre-boot configuration processing done.
Before=sysinit.target

[Install]
RequiredBy=sysinit.target
