# override/serial-getty@.service.d/override.conf
[Service]
EnvironmentFile=/etc/serial-getty/%i/serial-getty.env
ExecStart=
ExecStart=-/sbin/agetty $BAUD_RATE_LIST %I $TERM
